# metapkgs/guadalinexv8-desktop-conf/gcs/newfiles_skel/etc/init/welcome.conf
start on desktop-session-start
script
    if [ -e /etc/motd_welcome ]; then
        cat /etc/motd_welcome | /usr/bin/festival --tts
    fi
end script
