# metapkgs/guadalinexv8-desktop-conf/gcs/newfiles_skel/etc/init/welcome.conf
start on started lightdm
script
    if [ -e /etc/motd_welcome ]; then
        sleep 3;
        cat /etc/motd_welcome | /usr/bin/festival --tts
    fi
end script
